# docker_cpu/uwsgi/uwsgi.conf
start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

exec /usr/local/bin/uwsgi --emperor /home/dev/uwsgi

