# docker_cpu/ssh.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D
redirect_stderr=true
stdout_logfile=/var/log/ssh2.log
