# docker_cpu/Xvnc4.conf
[supervisord]
nodaemon=true

[program:vnc]
command=sh /etc/supervisor/conf.d/vncstart.sh
directory=/usr/bin/
redirect_stderr=true
stdout_logfile=/var/log/x11vnc.log 
